# root/home/.config/systemd/user/nas-sync-documents.service
[Unit]
Description=Sync Documents with NAS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-documents
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
