# root/home/.config/systemd/user/nas-sync-music.service
[Unit]
Description=Sync Music with NAS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-music
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
