# root/home/.config/systemd/user/openviking.service
[Unit]
Description=OpenViking Context Database
Documentation=https://github.com/volcengine/OpenViking
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/openviking-server
Restart=on-failure
RestartSec=5
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin"
Environment="OPENVIKING_CONFIG_FILE=%h/.openviking/ov.conf"
Environment="OPENVIKING_DATA_DIR=%h/.local/share/openviking"

[Install]
WantedBy=default.target