# root/home/.config/systemd/user/power-profile-switch.service
[Unit]
Description=Auto switch power profile based on AC state

[Service]
Type=oneshot
ExecStart=%h/.local/lib/power-profile-switch
