# root/etc/systemd/system/rkhunter-scan.timer
[Unit]
Description=Weekly rkhunter scan

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
