# home/.config/systemd/user/battery-notify.timer
[Unit]
Description=Battery notification timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min

[Install]
WantedBy=timers.target
