# home/.config/systemd/user/battery-notify.service
[Unit]
Description=Battery notification service

[Service]
Type=oneshot
ExecStart=%h/.local/lib/battery-notify
