# home/.config/systemd/user/nas-sync-audiobooks.service
[Unit]
Description=Sync Audiobooks with NAS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.local/bin/sync-audiobooks
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
